# crontab
* * * * * sudo sh /home/<user>/arp_check.sh
